# script/bypass7_2/bypass_monitor/silicom_multi.service
[Unit]
Description=Bypass STM

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/etc/init.d/bypass_silicom_enable.sh

[Install]
WantedBy=multi-user.target