# deployscripts/lcd/build/lcdd.service
[Unit]
Description=lcdd

[Service]
# EnvironmentFile=/etc/sysconfig/lm_sensors
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/lcdd /dev/ttyS1 &

[Install]
WantedBy=multi-user.target